# server/gps-web.service
[Unit]
Description=GPS Track Viewer Web UI
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/car-logger/server
EnvironmentFile=/opt/car-logger/server/.env
ExecStart=/opt/car-logger/server/.venv/bin/python -m gps_web.main
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
